# packages/tomfly-core/root/etc/init.d/tomfly
#!/bin/sh /etc/rc.common
# TomFly init script

START=99
STOP=10

PROG=/usr/bin/tomfly

start() {
    local enabled
    enabled=$(uci -q get tomfly.global.enabled 2>/dev/null)
    [ "${enabled:-1}" = "0" ] && return 0
    $PROG start
}

# Boot-time entry point (invoked by procd via the rc.d symlink). Gate auto-start
# on tomfly.global.autostart (default 0 = do NOT start on boot). Manual
# `/etc/init.d/tomfly start`, `tomfly start`, and the UI start button all call
# start() directly and are unaffected by this switch.
boot() {
    local autostart
    autostart=$(uci -q get tomfly.global.autostart 2>/dev/null)
    [ "${autostart:-0}" = "1" ] || return 0
    start
}

stop() {
    $PROG stop
}

reload() {
    $PROG reload
}
