# packages/luci-app-kanno/root/etc/uci-defaults/90-luci-app-kanno
#!/bin/sh
# First-boot init for KannoProxy LuCI app

# Reload uhttpd so it picks up new routes
/etc/init.d/uhttpd reload 2>/dev/null || true
/etc/init.d/rpcd reload 2>/dev/null || true
exit 0
